# miniconda/etc/profile.d/conda.csh
setenv CONDA_EXE "/mnt/2w12-data/audio-sampler-v2/miniconda/bin/conda";
setenv _CONDA_ROOT "/mnt/2w12-data/audio-sampler-v2/miniconda";
setenv _CONDA_EXE "/mnt/2w12-data/audio-sampler-v2/miniconda/bin/conda";
setenv CONDA_PYTHON_EXE "/mnt/2w12-data/audio-sampler-v2/miniconda/bin/python";
source "/mnt/2w12-data/audio-sampler-v2/miniconda/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh";
